# socat_xxx.ini
[program:socat_grafana]
command=/home/mbx/bin/socatuni.sh 47820 3000
autostart=yes
autorestart=yes
